# templates/i2p-router.service.j2
[Unit]
Description=I2P Router Service
After=network.target

[Service]
ExecStart={{ i2p_install_path }}/i2prouter console
User={{ i2p_user }}

[Install]
WantedBy=multi-user.target
